# tp3/aide-test.timer
[Unit]
Description=Timer AIDE

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
